# Faculdade/FIAP/Segundo_ano/1o_semestre/Computacao Neuromorfica/Checkpoint_1/lib/systemd/system/euclidean.service
[Unit]
Description=Euclidean Service
[Service]
Type=simple
ExecStart= /usr/bin/python3 /usr/bin/euclidean.py
Restart=on-failure
RestartSec=10
KillMode=process
[Install]
WantedBy=multi-user.target